# systemd/user/dualsense-persist.service
[Unit]
Description=Reapply persistent DualSense lightbar and player LED state
After=graphical-session.target

[Service]
Type=oneshot
# Same script as the udev hotplug path; exits 0 when no controller is bound.
ExecStart=/usr/local/bin/dualsense-init
