# systemd/user/dualsense-persist.timer
[Unit]
Description=Periodically reapply DualSense persistent colors

[Timer]
OnBootSec=15s
OnUnitActiveSec=30s
AccuracySec=1s
Unit=dualsense-persist.service

[Install]
WantedBy=timers.target
