# systemd/user/kopia-snapshot.timer
[Unit]
Description=Daily Kopia snapshot

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
